# i586con_brext/board/rootfs-overlay/etc/profile.d/01bashps1.sh
if [ -z "${BASH_VERSION-}" ] ; then
	# Shell is not bash, get out.
	return
fi

. /etc/bash/bashrc
