# brext/board/rootfs-overlay/etc/profile.d/02lang.sh
# Make all software aware that we're an UTF-8 system.
# (NB. musl has very limited NLS support, and we ship no translations.)
export LANG="en_US.UTF-8"
